# remarkable/tablet/quote-rotate.service
[Unit]
Description=Rotate reMarkable quote display image

[Service]
Type=oneshot
ExecStart=/home/root/quote-display/rotate.sh --skip-if-active
